# stage-tipi/01-config/files/tipi-runtipi-retry.service
[Unit]
Description=TipiOS — retry runTipi installation after failure
After=network-online.target docker.service
Wants=network-online.target
ConditionPathExists=/boot/firmware/tipi-install-failed.flag

[Service]
Type=oneshot
ExecStart=/opt/tipi-setup/retry-runtipi.sh
StandardOutput=journal
StandardError=journal
RemainAfterExit=no

[Install]
WantedBy=multi-user.target
