# stage-tipi/01-config/files/tipi-setup.service
[Unit]
Description=TipiOS — Premier démarrage : configuration et installation
Documentation=https://github.com/votre-repo/tipiOS
After=network-pre.target NetworkManager.service avahi-daemon.service NetworkManager-wait-online.service
Wants=NetworkManager.service NetworkManager-wait-online.service
ConditionPathExists=/var/lib/tipi-setup/.not-configured

[Service]
Type=simple
ExecStart=/opt/tipi-setup/start.sh
StandardOutput=journal
StandardError=journal
SyslogIdentifier=tipi-setup
# Redémarre uniquement si crash inattendu, pas si arrêt volontaire (code 0)
Restart=on-failure
RestartSec=10
StartLimitIntervalSec=120
StartLimitBurst=5

[Install]
WantedBy=multi-user.target
