# stage-tipi/01-config/files/tipi-setup.service
[Unit]
Description=RuntipiOS — Premier démarrage : configuration et installation
Documentation=https://github.com/votre-repo/RuntipiOS
After=network-pre.target avahi-daemon.service
Wants=avahi-daemon.service
ConditionPathExists=/var/lib/tipi-setup/.not-configured

[Service]
Type=simple
ExecStart=/opt/tipi-setup/start.sh
StandardOutput=append:/boot/firmware/tipi-setup.log
StandardError=append:/boot/firmware/tipi-setup.log
SyslogIdentifier=tipi-setup
# Redémarre uniquement si crash inattendu, pas si arrêt volontaire (code 0)
Restart=on-failure
RestartSec=10
StartLimitIntervalSec=120
StartLimitBurst=5

[Install]
WantedBy=multi-user.target
